# deploy/systemd/edgeway-uploader.service
[Unit]
Description=EdgeWay S3 yukleyici (tek tur)
[Service]
Type=oneshot
TimeoutStartSec=2h
User=edgeway
EnvironmentFile=/etc/edgeway/edgeway.env
WorkingDirectory=/opt/edgeway
ExecStart=/opt/edgeway/.venv/bin/python -m edgeway.uploader.s3sync
