# deploy/systemd/edgeway-recorder.service
[Unit]
Description=EdgeWay kayitci (ffmpeg segmenter)
After=network-online.target
[Service]
User=edgeway
EnvironmentFile=/etc/edgeway/edgeway.env
WorkingDirectory=/opt/edgeway
ExecStart=/opt/edgeway/.venv/bin/python -m edgeway.recorder.segmenter
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
